# nimbus/crontab
# crontab -l -u pi

# Backup home
20 4 * * * /home/pi/heapy-configuration/nimbus/backup.sh >> /tmp/cron.log 2>&1
